# packaging/app.bshk.agentvault.avd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Bundled LaunchAgent for SMAppService (macOS 13+). It is sealed inside the signed
  app at AgentVault.app/Contents/Library/LaunchAgents/ and registered by avd itself
  via SMAppService.agent(plistName:) (internal/loginitem/smappservice_darwin.m).

  BundleProgram is RELATIVE to the bundle root, so the signature stays valid wherever
  the .app is installed. Interactive ProcessType keeps avd in the user's GUI session
  so LocalAuthentication can present Touch ID. The build-from-source fallback uses a
  DIFFERENT plist with an absolute path — see internal/loginitem/launchagent.go.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>app.bshk.agentvault.avd</string>
    <key>BundleProgram</key>
    <string>Contents/MacOS/avd</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
